# crontab
0 * * * * /opt/venv/bin/python /opt/updater
